# build/rpm_specs/openstack-tacker-server.service
[Unit]
Description=OpenStack Tacker Server
After=syslog.target network.target

[Service]
Type=simple
User=tacker
ExecStart=/usr/bin/tacker-server --config-file /etc/tacker/tacker.conf --log-file /var/log/tacker/tacker.log
PrivateTmp=true
NotifyAccess=all
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
