# build/rpm_specs/openstack-congress-server.service
[Unit]
Description=OpenStack Congress Server
After=syslog.target network.target

[Service]
Type=simple
User=congress
ExecStart=/usr/bin/congress-server --config-file /etc/congress/congress.conf --log-file /var/log/congress/congress.log

[Install]
WantedBy=multi-user.target

